# crontab
SHELL=/bin/bash
BASH_ENV=/container.env

0 18 * * 7 /bin/bash /data/textpresso/tpctools/incremental_build.sh >& /tmp/incremental_build.log &

